# filesystem/HOME/.config/systemd/user/git-push@.timer
[Unit]
Description=Periodic git push for /%I

[Timer]
OnBootSec=15min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
