# filesystem/HOME/.config/systemd/user/git-shutdown-push@.service
[Unit]
Description=Final git push for /%I
DefaultDependencies=no
Before=shutdown.target

[Service]
Type=oneshot
WorkingDirectory=/%I
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/usr/bin/env nu --no-config-file --commands 'git add -A; try { git diff --cached --quiet; } except { git commit -m $$"[sync] (hostnamectl hostname)"; }; try { git push origin main; } except { notify-send "Shutdown push failed (%I)" -u critical; }'
